# .config/profile.d/chrome-flags.sh
export CHROME_FLAGS="--enable-features=WebUIDarkMode --force-dark-mode"
